# Marquee/marquee.service
[Unit]
Description=Marquee Display Daemon
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/marquee/marquee_daemon.py
WorkingDirectory=/home/marquee
User=marquee
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=SDL_VIDEODRIVER=KMSDRM
Environment=PYGAME_HIDE_SUPPORT_PROMPT=1
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=multi-user.target
